# service-files/restart.service
[Unit]
Description=Restart Arma service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart arma.service
